# deployment/notifications2.service
[Service]
WorkingDirectory=/opt/notifications
ExecStart=/usr/bin/node /opt/notifications/index.js
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=notifications-2
User=apiadmin
Group=apiadmin
Environment=NODE_ENV=production PORT=5010

[Install]
WantedBy=multi-user.target
